# docs/deploy/9router-mw.service
[Unit]
Description=9router-MW production gateway
Documentation=https://github.com/fazulfi/9router-mw
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=router
Group=router
WorkingDirectory=/opt/9router-mw/current
EnvironmentFile=/etc/9router-mw/env
# ipv4first via NODE_OPTIONS in EnvironmentFile
ExecStart=/usr/bin/node --max-old-space-size=4096 custom-server.js
Restart=always
RestartSec=3
TimeoutStopSec=30
KillSignal=SIGTERM
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/9router-mw /var/log/9router-mw
MemoryDenyWriteExecute=false
StandardOutput=journal
StandardError=journal
SyslogIdentifier=9router-mw

[Install]
WantedBy=multi-user.target
